# kite-trading-bot.service
[Unit]
Description=Kite Trading Bot with Kill Switch Monitoring
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/kite-algo
Environment="PATH=/home/ubuntu/kite-algo/.venv/bin"
ExecStart=/home/ubuntu/kite-algo/.venv/bin/python /home/ubuntu/kite-algo/start_bot_with_monitor.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/kite-algo/logs/bot_service.log
StandardError=append:/home/ubuntu/kite-algo/logs/bot_service_error.log

[Install]
WantedBy=multi-user.target
