# deploy/kite-algo-celery-beat.service
[Unit]
Description=Kite Algo Trading Platform - Celery Beat Scheduler
After=network.target redis-server.service
Wants=redis-server.service kite-algo-celery.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/kite-algo
EnvironmentFile=/var/www/kite-algo/.env
ExecStart=/var/www/kite-algo/venv/bin/celery -A src.workers.celery_app beat \
    --loglevel=info \
    --logfile=/var/log/kite-algo/celery-beat.log
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
